# dashboard-site/com.dashboard.gunicorn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 服務識別碼 -->
    <key>Label</key>
    <string>com.dashboard.gunicorn</string>
    
    <!-- 啟動程式 -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>gunicorn</string>
        <string>-c</string>
        <string>/Users/aiserver/.openclaw/workspace/dashboard-site/gunicorn.conf.py</string>
        <string>-w</string>
        <string>6</string>
        <string>-b</string>
        <string>127.0.0.1:3000</string>
        <string>--timeout</string>
        <string>60</string>
        <string>--access-logfile</string>
        <string>/Users/aiserver/srv/logs/gunicorn_access.log</string>
        <string>--error-logfile</string>
        <string>/Users/aiserver/srv/logs/gunicorn_error.log</string>
        <string>app:app</string>
    </array>
    
    <!-- 工作目錄 -->
    <key>WorkingDirectory</key>
    <string>/Users/aiserver/.openclaw/workspace/dashboard-site</string>
    
    <!-- 環境變數 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/aiserver/.openclaw/workspace/dashboard-site</string>
        <key>DB_PATH</key>
        <string>/Users/aiserver/srv/db/company.db</string>
    </dict>
    
    <!-- 自動重啟設定 -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- 當機後重啟延遲（秒） -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- 標準輸出/錯誤 -->
    <key>StandardOutPath</key>
    <string>/Users/aiserver/srv/logs/gunicorn_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aiserver/srv/logs/gunicorn_launchd_error.log</string>
    
    <!-- 背景執行 -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- 使用者 -->
    <key>UserName</key>
    <string>aiserver</string>
</dict>
</plist>
